# scripts/wifi-monitor.service.template
[Unit]
Description=Wi-Fi Monitor Dashboard
After=network.target

[Service]
Type=simple
User={USER}
WorkingDirectory={PROJECT_DIR}
Environment="PATH={PROJECT_DIR}/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Run with 1 worker and 4 threads for the background auto-ping loop and SSE stream
ExecStart={PROJECT_DIR}/venv/bin/gunicorn -w 1 --threads 4 -b 0.0.0.0:5000 "dashboard.app:create_app(auto_run=True)"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
